# templates/podman-run@CONTAINER.container
[Container]
[Container]
Image={{ container.image }}
