# files/podman-run@.container
[Unit]
Description=Podman single-container Quadlet template

[Service]
Restart=always

[Install]
WantedBy=default.target
